# src/server/supervisord.conf
[program:node]
command = nodemon server.coffee
directory = /home/dotcloud/current
