# src/server/supervisord.conf
[program:node]
environment = NODE_ENV=production
command = /home/dotcloud/current/node_modules/.bin/coffee ./src/server/server.coffee
directory = /home/dotcloud/current
