# src/server/supervisord.conf
[program:node]
command = node server.coffee
directory = /home/dotcloud/current
