# supervisord.conf
[program:node]
command = /home/dotcloud/current/node_modules/.bin/coffee src/server/server.coffee
directory = /home/dotcloud/current
